# scripts/dashboard-server.service
[Unit]
Description=HTTP Server Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/etc/dashboard-server
ExecStart=/usr/local/bin/dashboard-server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
